# src/tools/akamai/osd_debug/etc/cron.d/osd_debug
*/5 * * * * /srv/exporter_scripts/socket_stat.sh
